# config/packages/billgastro-src/etc/cron.monthly/1billgastro
#!/bin/sh
:>/opt/billgastro/billgastro/log/production.log
