# config/debian-package/etc/profile.d/salor_hospitality.sh
#!/bin/sh
export RAILS_ENV=production
export PATH=$PATH:/opt/salor/ruby/bin
