# debian/salor-hospitality.cron.d
#
# Regular cron jobs for the salor-hospitality package
#
0 4 * * * /usr/bin/salor-maintainance h
